# .config/systemd/user/battery-warning.service
[Unit]
Description=Battery Low Warning

[Service]
Type=oneshot
ExecStart=/home/mentat/.battery/battery-warning.sh
