# .config/systemd/user/check-updates.service
[Unit]
Description=Notify of available system updates

[Service]
Type=oneshot
ExecStart=/home/mentat/.check-updates.sh
